# utils/supervisord.conf.d/application.conf
[program:application]
directory=%(ENV_DIRECTORY)s
command=wine explorer /desktop=application,%(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)s %(ENV_APPLICATION)s %(ENV_ARGUMENTS)s
autorestart=true
