# utils/supervisord.conf.d/application.conf
[program:application]
directory=%(ENV_WORKDIR)s
command=/root/utils/application.sh
autorestart=true
